# deploy/jailtimed.service
[Unit]
Description=Jailtime daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/sbin/jailtimed --config /etc/jailtime/jail.yaml
Restart=on-failure
RestartSec=2s
RuntimeDirectory=jailtime
RuntimeDirectoryMode=0755

# Hardening (tune for firewall backend needs)
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
